# deploy/all-in-one/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:postgres]
command=/bin/bash -lc 'POSTGRES_BIN="$(find /usr -type f -name postgres 2>/dev/null | head -n1)"; exec "$POSTGRES_BIN" -D /var/lib/postgresql/data'
user=postgres
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/postgres.out.log
stderr_logfile=/var/log/supervisor/postgres.err.log
priority=10

[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf --daemonize no
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/redis.out.log
stderr_logfile=/var/log/supervisor/redis.err.log
priority=20

[program:minio]
command=/usr/local/bin/minio server /data/minio --console-address :9001
directory=/
autostart=true
autorestart=true
environment=MINIO_ROOT_USER="%(ENV_MINIO_ROOT_USER)s",MINIO_ROOT_PASSWORD="%(ENV_MINIO_ROOT_PASSWORD)s"
stdout_logfile=/var/log/supervisor/minio.out.log
stderr_logfile=/var/log/supervisor/minio.err.log
priority=30

[program:keycloak]
command=/opt/keycloak/bin/kc.sh start --http-enabled=true --hostname-strict=false --proxy-headers=xforwarded --http-management-port=9002 --http-management-host=127.0.0.1 --import-realm
autostart=true
autorestart=true
environment=KC_BOOTSTRAP_ADMIN_USERNAME="%(ENV_KEYCLOAK_ADMIN_USERNAME)s",KC_BOOTSTRAP_ADMIN_PASSWORD="%(ENV_KEYCLOAK_ADMIN_PASSWORD)s",KC_DB="postgres",KC_DB_URL="jdbc:postgresql://127.0.0.1:5432/%(ENV_KC_POSTGRES_DB)s",KC_DB_USERNAME="%(ENV_POSTGRES_USER)s",KC_DB_PASSWORD="%(ENV_POSTGRES_PASSWORD)s",KC_HEALTH_ENABLED="true",KC_METRICS_ENABLED="true",KC_HTTP_MANAGEMENT_PORT="9002",KC_HTTP_MANAGEMENT_HOST="127.0.0.1"
stdout_logfile=/var/log/supervisor/keycloak.out.log
stderr_logfile=/var/log/supervisor/keycloak.err.log
priority=40

[program:backend]
command=/usr/local/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --log-level info
directory=/code
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/backend.out.log
stderr_logfile=/var/log/supervisor/backend.err.log
priority=50

[program:celery_worker]
command=/usr/local/bin/celery -A app.celery_app worker --loglevel=info --concurrency=1 --queues=default,processing,reports,ai_llm -n worker@%%h
directory=/code
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/celery_worker.out.log
stderr_logfile=/var/log/supervisor/celery_worker.err.log
priority=60

[program:celery_beat]
command=/usr/local/bin/celery -A app.celery_app beat --loglevel=info
directory=/code
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/celery_beat.out.log
stderr_logfile=/var/log/supervisor/celery_beat.err.log
priority=70

[program:flower]
command=/usr/local/bin/celery -A app.celery_app flower --port=5555 --broker_api=redis://127.0.0.1:6379/0 --basic-auth=%(ENV_FLOWER_BASIC_AUTH)s
directory=/code
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/flower.out.log
stderr_logfile=/var/log/supervisor/flower.err.log
priority=80
